# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
cd /config || exit
exec \
	nginx -c /etc/nginx/nginx.conf -g "daemon off;"